# crontab
*/5 * * * * /opt/services/coronacircles/manage.py mail_participants